# etc/service/portmanager/run
#!/bin/sh -e

cd /data

export PLUGIN_PORT=9001

exec /usr/local/bin/portmanager 2>&1
