# etc/service/frps/run
#!/bin/sh -e

/usr/local/bin/gotemp -f /etc/frps.ini.tpl -o /etc/frps.ini

exec /usr/local/bin/frps -c /etc/frps.ini 2>&1
